# GoldenDict/home_config_autostart/GoldenDict.desktop
[Desktop Entry]
Type=Application
Name=GoldenDict
Icon=goldendict
Comment=
Exec=goldendict
StartupNotify=true
Terminal=false
Hidden=false

